# scripts/daemon/com.memxp.sync-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.memxp.sync-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/memxp-daemon-guard</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <!-- Only restart if guard exited non-zero (daemon died) -->
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Prevent rapid restart loops: 60s cooldown between attempts -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- Allow 30s for graceful shutdown (daemon stop) before SIGKILL -->
    <key>ExitTimeOut</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/.memxp/logs/daemon-guard.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.memxp/logs/daemon-guard.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
