# scripts/daemon/memxp-daemon.service
[Unit]
Description=memxp sync daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/root/.vaultp2p/env.systemd
EnvironmentFile=-/root/.memxp/env.systemd
ExecStart=/usr/local/bin/memxp-daemon-guard
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
